# MacScope-vHID-Relay/com.macscope.vhid.relay.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.macscope.vhid.relay</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>
            LOG="/var/log/macscope-vhid.log"
            echo "[startup] Waiting for Karabiner VirtualHID driver..." >> "$LOG"
            for i in {1..15}; do
                if systemextensionsctl list | grep -q "Karabiner-DriverKit-VirtualHIDDevice.*activated.*enabled"; then
                    if pgrep -x "Karabiner-VirtualHIDDevice-Daemon" >/dev/null 2>&1; then
                        echo "[ok] Karabiner VirtualHID ready after $i attempts" >> "$LOG"
                        break
                    fi
                fi
                echo "[wait] Attempt $i - driver not ready" >> "$LOG"
                sleep 2
            done
            mkdir -p /var/run/macscope
            chmod 777 /var/run/macscope
            echo "[launch] Starting macscope-vhid relay" >> "$LOG"
            exec /usr/local/bin/macscope-vhid --verbose >> "$LOG" 2>&1
        </string>
    </array>

    <key>RunAtLoad</key><true/>
    <key>KeepAlive</key><true/>
    <key>StandardOutPath</key><string>/var/log/macscope-vhid.log</string>
    <key>StandardErrorPath</key><string>/var/log/macscope-vhid.log</string>
</dict>
</plist>
